# com.malko.meridian-api.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.malko.meridian-api</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/malkobot/.openclaw/workspace/projects/meridian-api/venv/bin/uvicorn</string>
        <string>app.main:app</string>
        <string>--host</string>
        <string>0.0.0.0</string>
        <string>--port</string>
        <string>8100</string>
        <string>--workers</string>
        <string>4</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/malkobot/.openclaw/workspace/projects/meridian-api</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/malkobot</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/meridian-api.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/meridian-api.error.log</string>

    <key>SoftResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>65536</integer>
    </dict>
    <key>HardResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>65536</integer>
    </dict>
</dict>
</plist>
